# services/stiximaging-plotter.service
[Unit]
Description=STIX data processing pipeline

[Service]
User=root
ExecStart=/opt/stix/pystix/bin/plot_idl
Restart=always

[Install]
WantedBy=default.target
